# evaluation 2.0/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:dashboard]
command=python -m uvicorn app:app --host 127.0.0.1 --port 9010 --proxy-headers --forwarded-allow-ips=*
directory=/app/dashboard
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:fastapi-tx]
command=python -m uvicorn app.main:app --host 127.0.0.1 --port 8000 --root-path /services/tx
directory=/app/fastapi_transactions
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:fraud-score]
command=python -m uvicorn app.main:app --host 127.0.0.1 --port 8001 --root-path /services/fraud
directory=/app/fraud-score-system/service
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:fintech-demo]
command=python -m uvicorn app.main:app --host 127.0.0.1 --port 8002 --root-path /services/fintech
directory=/app/sample-projects/fintech-platform
environment=PYTHONPATH="/app/sample-projects/fintech-platform"
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
